# plugin/codex-watcher/com.codzr.codex-completion-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codzr.codex-completion-watcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/codzr/.config/opencode/plugin/codex-watcher/run-codex-completion-watcher.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OPENCODE_NOTIFY_TOKEN</key>
    <string>__REPLACE_WITH_LOCAL_OPENCODE_NOTIFY_TOKEN__</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/codzr/Library/Logs/codex-completion-watcher.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/codzr/Library/Logs/codex-completion-watcher.error.log</string>
</dict>
</plist>
